# postgres/rootfs/etc/services.d/postgresql/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start PostgreSQL service
# ==============================================================================

if ! bashio::fs.directory_exists "${POSTGRES_DATA}"; then
    bashio::log.info "Create database directory \"${POSTGRES_DATA}\""
    mkdir "${POSTGRES_DATA}"
    chown postgres:postgres "${POSTGRES_DATA}"

    exec s6-setuidgid postgres /usr/libexec/postgres/initdb.sh
else
    bashio::log.info "Using existing postgres initial system"
    exec s6-setuidgid postgres /usr/libexec/postgres/checkdb.sh
fi

# Execution follows following steps:
#  initdb OR checkdb
#  start
#  config
#  monitor
