# backup-s3/rootfs/etc/services.d/backup/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Setting AWS CLI options."

BUCKET="$(bashio::config 'bucketname')"
ENDPOINT="$(bashio::config 's3endpoint')"
KEY="$(bashio::config 's3keyid')"
SECRET="$(bashio::config 's3secret')"

aws configure set aws_access_key_id "${KEY}"
aws configure set aws_secret_access_key "${SECRET}"
aws configure set plugins.endpoint awscli_plugin_endpoint

bashio::log.info "Starting backup synchronisation to S3 bucket ${BUCKET} on endpoint ${ENDPOINT}."

if aws s3 sync /backup "s3://${BUCKET}/" --delete --endpoint-url "${ENDPOINT}"
then
    bashio::log.info "Synchronization succeeded."
    exit 0
else
    bashio::log.warning "Synchronization failed."
    exit 1
fi
